# profile.d/aws.sh
export AWS_DEFAULT_PROFILE=cassandra-benchmark
export AWS_DEFAULT_REGION=us-west-2
